# needlesslaser.service
[Unit]
Description=needlesslaser cloud9 initialization

[Service]
ExecStart=/opt/needlesslaser/sbin/run

[Install]
WantedBy=multi-user.target